# software/S0192/S0192-enterprise-layer.json
{"description": "Enterprise techniques used by Pupy, ATT&CK software S0192 v1.2", "name": "Pupy (S0192)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"techniqueID": "T1548", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1548.002", "showSubtechniques": true, "comment": "Pupy can bypass Windows UAC through either DLL hijacking, eventvwr, or appPaths.[1]"}, {"techniqueID": "T1134", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1134.001", "showSubtechniques": true, "comment": "Pupy can obtain a list of SIDs and provide the option for selecting process tokens to impersonate.[1]"}, {"techniqueID": "T1087", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1087.001", "showSubtechniques": true, "comment": "Pupy uses PowerView and Pywerview to perform discovery commands such as net user, net group, net local group, etc.[1]"}, {"techniqueID": "T1071", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1071.001", "showSubtechniques": true, "comment": "Pupy can communicate over HTTP for C2.[1]"}, {"techniqueID": "T1560", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1560.001", "showSubtechniques": true, "comment": "Pupy can compress data with Zip before sending it over C2.[1]"}, {"score": 1, "techniqueID": "T1123", "showSubtechniques": false, "comment": "Pupy can record sound with the microphone.[1]"}, {"techniqueID": "T1547", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1547.001", "showSubtechniques": true, "comment": "Pupy adds itself to the startup folder or adds itself to the Registry key SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run for persistence.[1]"}, {"techniqueID": "T1059", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1059.001", "showSubtechniques": true, "comment": "Pupy has a module for loading and executing PowerShell scripts.[1]"}, {"score": 1, "techniqueID": "T1059.006", "showSubtechniques": true, "comment": "Pupy can use an add on feature when creating payloads that allows you to create custom Python scripts (\"scriptlets\") to perform tasks offline (without requiring a session) such as sandbox detection, adding persistence, etc.[1]"}, {"techniqueID": "T1136", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1136.001", "showSubtechniques": true, "comment": "Pupy can user PowerView to execute \"net user\" commands and create local system accounts.[1]"}, {"score": 1, "techniqueID": "T1136.002", "showSubtechniques": true, "comment": "Pupy can user PowerView to execute \"net user\" commands and create domain accounts.[1]"}, {"techniqueID": "T1543", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1543.002", "showSubtechniques": true, "comment": "Pupy can be used to establish persistence using a systemd service.[1]"}, {"score": 1, "techniqueID": "T1555", "showSubtechniques": true, "comment": "Pupy can use Lazagne for harvesting credentials.[1]"}, {"score": 1, "techniqueID": "T1555.003", "showSubtechniques": true, "comment": "Pupy can use Lazagne for harvesting credentials.[1]"}, {"techniqueID": "T1114", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1114.001", "showSubtechniques": true, "comment": "Pupy can interact with a victim\u2019s Outlook session and look through folders and emails.[1]"}, {"score": 1, "techniqueID": "T1573", "showSubtechniques": true, "comment": "Pupy can use Obfs3, a pluggable transport, to add another layer of encryption and obfuscate TLS.[1]"}, {"score": 1, "techniqueID": "T1573.002", "showSubtechniques": true, "comment": "Pupy's default encryption for its C2 communication channel is SSL, but it also has transport options for RSA and AES.[1]"}, {"score": 1, "techniqueID": "T1041", "showSubtechniques": false, "comment": "Pupy can send screenshots files, keylogger data, files, and recorded audio back to the C2 server.[1]"}, {"score": 1, "techniqueID": "T1083", "showSubtechniques": false, "comment": "Pupy can walk through directories and recursively search for strings in files.[1]"}, {"techniqueID": "T1551", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1551.001", "showSubtechniques": true, "comment": "Pupy has a module to clear event logs with PowerShell.[1]"}, {"score": 1, "techniqueID": "T1105", "showSubtechniques": false, "comment": "Pupy can upload and download to/from a victim machine.[1]"}, {"techniqueID": "T1056", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1056.001", "showSubtechniques": true, "comment": "Pupy uses a keylogger to capture keystrokes it then sends back to the server after it is stopped.[1]"}, {"techniqueID": "T1557", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1557.001", "showSubtechniques": true, "comment": "Pupy can sniff plaintext network credentials and use NBNS Spoofing to poison name services.[1]"}, {"score": 1, "techniqueID": "T1046", "showSubtechniques": false, "comment": "Pupy has a built-in module for port scanning.[1]"}, {"score": 1, "techniqueID": "T1135", "showSubtechniques": false, "comment": "Pupy can list local and remote shared drives and folders over SMB.[1]"}, {"techniqueID": "T1003", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1003.001", "showSubtechniques": true, "comment": "Pupy can execute Lazagne as well as Mimikatz using PowerShell.[1]"}, {"score": 1, "techniqueID": "T1003.004", "showSubtechniques": true, "comment": "Pupy can use Lazagne for harvesting credentials.[1]"}, {"score": 1, "techniqueID": "T1003.005", "showSubtechniques": true, "comment": "Pupy can use Lazagne for harvesting credentials.[1]"}, {"score": 1, "techniqueID": "T1057", "showSubtechniques": false, "comment": "Pupy can list the running processes and get the process ID and parent process\u2019s ID.[1]"}, {"techniqueID": "T1055", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1055.001", "showSubtechniques": true, "comment": "Pupy can migrate into another process using reflective DLL injection.[1]"}, {"techniqueID": "T1021", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1021.001", "showSubtechniques": true, "comment": "Pupy can enable/disable RDP connection and can start a remote desktop session using a browser web socket client.[1]"}, {"score": 1, "techniqueID": "T1113", "showSubtechniques": false, "comment": "Pupy can drop a mouse-logger that will take small screenshots around at each click and then send back to the server.[1]"}, {"score": 1, "techniqueID": "T1082", "showSubtechniques": false, "comment": "Pupy can grab a system\u2019s information including the OS version, architecture, etc.[1]"}, {"score": 1, "techniqueID": "T1016", "showSubtechniques": false, "comment": "Pupy has built in commands to identify a host\u2019s IP address and find out other network configuration settings by viewing connected sessions.[1]"}, {"score": 1, "techniqueID": "T1049", "showSubtechniques": false, "comment": "Pupy has a built-in utility command for netstat, can do net session through PowerView, and has an interactive shell which can be used to discover additional information.[1]"}, {"score": 1, "techniqueID": "T1033", "showSubtechniques": false, "comment": "Pupy can enumerate local information for Linux hosts and find currently logged on users for Windows hosts.[1]"}, {"techniqueID": "T1569", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1569.002", "showSubtechniques": true, "comment": "Pupy uses PsExec to execute a payload or commands on a remote host.[1]"}, {"techniqueID": "T1552", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1552.001", "showSubtechniques": true, "comment": "Pupy can use Lazagne for harvesting credentials.[1]"}, {"techniqueID": "T1550", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1550.003", "showSubtechniques": true, "comment": "Pupy can also perform pass-the-ticket.[1]"}, {"score": 1, "techniqueID": "T1125", "showSubtechniques": false, "comment": "Pupy can access a connected webcam and capture pictures.[1]"}, {"techniqueID": "T1497", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1497.001", "showSubtechniques": true, "comment": "Pupy has a module that checks a number of indicators on the system to determine if its running on a virtual machine.[1]"}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by Pupy", "color": "#66b1ff"}]}