# software/S0513/S0513-enterprise-layer.json
{"description": "Enterprise techniques used by LiteDuke, ATT&CK software S0513 (v1.0)", "name": "LiteDuke (S0513)", "domain": "enterprise-attack", "versions": {"layer": "4.5", "attack": "17", "navigator": "5.1.0"}, "techniques": [{"techniqueID": "T1071", "showSubtechniques": true}, {"techniqueID": "T1071.001", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) can use HTTP GET requests in C2 communications.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1547", "showSubtechniques": true}, {"techniqueID": "T1547.001", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) can create persistence by adding a shortcut in the CurrentVersion\\Run Registry key.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1140", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) has the ability to decrypt and decode multiple layers of obfuscation.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1070", "showSubtechniques": true}, {"techniqueID": "T1070.004", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) can securely delete files by first writing random data to the file.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1105", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) has the ability to download files.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1027", "showSubtechniques": true}, {"techniqueID": "T1027.002", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) has been packed with multiple layers of encryption.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1027.003", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) has used image files to hide its loader component.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1012", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) can query the Registry to check for the presence of HKCU\\Software\\KasperskyLab.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1518", "showSubtechniques": true}, {"techniqueID": "T1518.001", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) has the ability to check for the presence of Kaspersky security software.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1082", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) can enumerate the CPUID and BIOS version on a compromised system.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1016", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) has the ability to discover the proxy configuration of Firefox and/or Opera.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1033", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) can enumerate the account name on a targeted system.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1497", "showSubtechniques": true}, {"techniqueID": "T1497.003", "comment": "[LiteDuke](https://attack.mitre.org/software/S0513) can wait 30 seconds before executing additional code if security software is detected.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by LiteDuke", "color": "#66b1ff"}]}