# versions/v7/software/S0062/S0062-enterprise-layer.json
{"description": "Enterprise techniques used by DustySky, ATT&CK software S0062 v1.1", "name": "DustySky (S0062)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"techniqueID": "T1071", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1071.001", "showSubtechniques": true, "comment": "DustySky has used both HTTP and HTTPS for C2.[1]"}, {"techniqueID": "T1560", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1560.001", "showSubtechniques": true, "comment": "DustySky can compress files via RAR while staging data to be exfiltrated.[3]"}, {"techniqueID": "T1547", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1547.001", "showSubtechniques": true, "comment": "DustySky achieves persistence by creating a Registry entry in HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run.[1]"}, {"techniqueID": "T1074", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1074.001", "showSubtechniques": true, "comment": "DustySky created folders in temp directories to host collected files before exfiltration.[3]"}, {"score": 1, "techniqueID": "T1041", "showSubtechniques": false, "comment": "DustySky has exfiltrated data to the C2 server.[3]"}, {"score": 1, "techniqueID": "T1008", "showSubtechniques": false, "comment": "DustySky has two hard-coded domains for C2 servers; if the first does not respond, it will try the second.[1]"}, {"score": 1, "techniqueID": "T1083", "showSubtechniques": false, "comment": "DustySky scans the victim for files that contain certain keywords and document types including PDF, DOC, DOCX, XLS, and XLSX, from a list that is obtained from the C2 as a text file. It can also identify logical drives for the infected machine.[1][3]"}, {"techniqueID": "T1070", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1070.004", "showSubtechniques": true, "comment": "DustySky can delete files it creates from the infected system.[3]"}, {"techniqueID": "T1056", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1056.001", "showSubtechniques": true, "comment": "DustySky contains a keylogger.[1]"}, {"score": 1, "techniqueID": "T1570", "showSubtechniques": false, "comment": "DustySky searches for network drives and removable media and duplicates itself onto them.[1]"}, {"score": 1, "techniqueID": "T1027", "showSubtechniques": false, "comment": "The DustySky dropper uses a function to obfuscate the name of functions and other parts of the malware.[1]"}, {"score": 1, "techniqueID": "T1120", "showSubtechniques": false, "comment": "DustySky can detect connected USB devices.[3]"}, {"score": 1, "techniqueID": "T1057", "showSubtechniques": false, "comment": "DustySky collects information about running processes from victims.[1][3]"}, {"score": 1, "techniqueID": "T1091", "showSubtechniques": false, "comment": "DustySky searches for removable media and duplicates itself onto it.[1]"}, {"score": 1, "techniqueID": "T1113", "showSubtechniques": false, "comment": "DustySky captures PNG screenshots of the main screen.[3]"}, {"score": 1, "techniqueID": "T1518", "showSubtechniques": true, "comment": "DustySky lists all installed software for the infected machine.[3]"}, {"score": 1, "techniqueID": "T1518.001", "showSubtechniques": true, "comment": "DustySky checks for the existence of anti-virus.[1]"}, {"score": 1, "techniqueID": "T1082", "showSubtechniques": false, "comment": "DustySky extracts basic information about the operating system.[1]"}, {"score": 1, "techniqueID": "T1529", "showSubtechniques": false, "comment": "DustySky can shutdown the infected machine.[3]"}, {"score": 1, "techniqueID": "T1047", "showSubtechniques": false, "comment": "The DustySky dropper uses Windows Management Instrumentation to extract information about the operating system and whether an anti-virus is active.[1]"}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by DustySky", "color": "#66b1ff"}]}