# versions/v7/groups/G0087/G0087-enterprise-layer.json
{"description": "Enterprise techniques used by APT39, ATT&CK group G0087 v2.1", "name": "APT39 (G0087)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"techniqueID": "T1071", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1071.004", "showSubtechniques": true, "comment": "APT39 has used remote access tools that leverage DNS in communications with C2.[4]"}, {"score": 1, "techniqueID": "T1071.001", "showSubtechniques": true, "comment": "APT39 has used HTTP in communications with C2.[4]"}, {"techniqueID": "T1560", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1560.001", "showSubtechniques": true, "comment": "APT39 has used WinRAR and 7-Zip to compress an archive stolen data. [1]"}, {"techniqueID": "T1547", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1547.001", "showSubtechniques": true, "comment": "APT39 has maintained persistence using the startup folder. [1]"}, {"score": 1, "techniqueID": "T1547.009", "showSubtechniques": true, "comment": "APT39 has modified LNK shortcuts. [1]"}, {"score": 1, "techniqueID": "T1110", "showSubtechniques": false, "comment": "APT39 has used Ncrack to reveal credentials.[1]"}, {"score": 1, "techniqueID": "T1115", "showSubtechniques": false, "comment": "APT39 has used tools capable of stealing contents of the clipboard.[3]"}, {"techniqueID": "T1059", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1059.003", "showSubtechniques": true, "comment": "APT39 utilized custom scripts to perform internal reconnaissance. [1]"}, {"score": 1, "techniqueID": "T1059.001", "showSubtechniques": true, "comment": "APT39 has used PowerShell to execute malicious code.[4][3]"}, {"score": 1, "techniqueID": "T1059.006", "showSubtechniques": true, "comment": "APT39 has used a command line utility and a network scanner written in python.[4]"}, {"techniqueID": "T1136", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1136.001", "showSubtechniques": true, "comment": "APT39 has created accounts on multiple compromised hosts to perform actions within the network.[4]"}, {"score": 1, "techniqueID": "T1555", "showSubtechniques": false, "comment": "APT39 has used the Smartftp Password Decryptor tool to decrypt FTP passwords.[4]"}, {"score": 1, "techniqueID": "T1005", "showSubtechniques": false, "comment": "APT39 has used a tool to steal files from the compromised host.[3]"}, {"score": 1, "techniqueID": "T1190", "showSubtechniques": false, "comment": "APT39 has used SQL injection for initial compromise.[3]"}, {"score": 1, "techniqueID": "T1105", "showSubtechniques": false, "comment": "APT39 has downloaded tools to compromised hosts.[3]"}, {"techniqueID": "T1056", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1056.001", "showSubtechniques": true, "comment": "APT39 has used tools for capturing keystrokes.[3]"}, {"techniqueID": "T1036", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1036.005", "showSubtechniques": true, "comment": "APT39 has used a tool named mfevtpse.exe to proxy C2 communications, closely mimicking a legitimate McAfee file mfevtps.exe.[4]"}, {"score": 1, "techniqueID": "T1046", "showSubtechniques": false, "comment": "APT39 has used CrackMapExec and a custom port scanner known as BLUETORCH for network scanning [1][4]"}, {"score": 1, "techniqueID": "T1135", "showSubtechniques": false, "comment": "APT39 has used the post exploitation tool CrackMapExec to enumerate network shares.[4]"}, {"techniqueID": "T1027", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1027.002", "showSubtechniques": true, "comment": "APT39 has packed tools with UPX, and has repacked a modified version of Mimikatz to thwart anti-virus detection.[1][4]"}, {"score": 1, "techniqueID": "T1003", "showSubtechniques": true, "comment": "APT39 has used different versions of Mimikatz to obtain credentials.[4]"}, {"score": 1, "techniqueID": "T1003.001", "showSubtechniques": true, "comment": "APT39 has used Mimikatz, Windows Credential Editor and ProcDump to dump credentials. [1]"}, {"techniqueID": "T1566", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1566.001", "showSubtechniques": true, "comment": "APT39 leveraged spearphishing emails with malicious attachments to initially compromise victims. [1][3]"}, {"score": 1, "techniqueID": "T1566.002", "showSubtechniques": true, "comment": "APT39 leveraged spearphishing emails with malicious links to initially compromise victims. [1]"}, {"techniqueID": "T1090", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1090.001", "showSubtechniques": true, "comment": "APT39 used custom tools to create SOCK5 and custom protocol proxies between infected hosts. [1][4]"}, {"score": 1, "techniqueID": "T1090.002", "showSubtechniques": true, "comment": "APT39 has used various tools to proxy C2 communications.[4]"}, {"techniqueID": "T1021", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1021.001", "showSubtechniques": true, "comment": "APT39 has been seen using RDP for lateral movement and persistence, in some cases employing the rdpwinst tool for mangement of multiple sessions. [1][4]"}, {"score": 1, "techniqueID": "T1021.004", "showSubtechniques": true, "comment": "APT39 used secure shell (SSH) to move laterally among their targets. [1]"}, {"score": 1, "techniqueID": "T1021.002", "showSubtechniques": true, "comment": "APT39 has used SMB for lateral movement.[3]"}, {"score": 1, "techniqueID": "T1018", "showSubtechniques": false, "comment": "APT39 has used nbtscan and custom tools to discover remote systems. [1][4][3]"}, {"techniqueID": "T1053", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1053.005", "showSubtechniques": true, "comment": "APT39 has created scheduled tasks for persistence. [1][4]"}, {"score": 1, "techniqueID": "T1113", "showSubtechniques": false, "comment": "APT39 has used a screen capture utility to take screenshots on a compromised host.[3]"}, {"techniqueID": "T1505", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1505.003", "showSubtechniques": true, "comment": "APT39 has installed ANTAK and ASPXSPY web shells. [1]"}, {"score": 1, "techniqueID": "T1033", "showSubtechniques": false, "comment": "APT39 used Remexi to collect usernames from the system. [2]"}, {"techniqueID": "T1569", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1569.002", "showSubtechniques": true, "comment": "APT39 has used post-exploitation tools including RemCom and the Non-sucking Service Manager (NSSM) to execute processes.[4][3]"}, {"techniqueID": "T1204", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1204.002", "showSubtechniques": true, "comment": "APT39 has sent spearphishing emails in an attempt to lure users to click on a malicious attachment. [1][4][3]"}, {"score": 1, "techniqueID": "T1204.001", "showSubtechniques": true, "comment": "APT39 has sent spearphishing emails in an attempt to lure users to click on a malicious link. [1]"}, {"score": 1, "techniqueID": "T1078", "showSubtechniques": false, "comment": "APT39 has used stolen credentials to compromise Outlook Web Access (OWA). [1]"}, {"techniqueID": "T1102", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1102.002", "showSubtechniques": true, "comment": "APT39 has communicated with C2 through files uploaded to DropBox.[4]"}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by APT39", "color": "#66b1ff"}]}