# software/S0369/S0369-enterprise-layer.json
{"description": "Enterprise techniques used by CoinTicker, ATT&CK software S0369 (v1.1)", "name": "CoinTicker (S0369)", "domain": "enterprise-attack", "versions": {"layer": "4.5", "attack": "17", "navigator": "5.1.0"}, "techniques": [{"techniqueID": "T1059", "showSubtechniques": true}, {"techniqueID": "T1059.003", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) executes a bash script to establish a reverse shell.(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1059.004", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) executes a bash script to establish a reverse shell.(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1059.006", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) executes a Python script to download its second stage.(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1543", "showSubtechniques": true}, {"techniqueID": "T1543.001", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) creates user launch agents named .espl.plist and com.apple.[random string].plist to establish persistence.(Citation: CoinTicker 2019)\t", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1140", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) decodes the initially-downloaded hidden encoded file using OpenSSL.(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1564", "showSubtechniques": true}, {"techniqueID": "T1564.001", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) downloads the following hidden files to evade detection and maintain persistence: /private/tmp/.info.enc, /private/tmp/.info.py, /private/tmp/.server.sh, ~/Library/LaunchAgents/.espl.plist, ~/Library/Containers/.[random string]/[random string].(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1105", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) executes a Python script to download its second stage.(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1027", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) initially downloads a hidden encoded file.(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1553", "showSubtechniques": true}, {"techniqueID": "T1553.001", "comment": "[CoinTicker](https://attack.mitre.org/software/S0369) downloads the EggShell mach-o binary using curl, which does not set the quarantine flag.(Citation: CoinTicker 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by CoinTicker", "color": "#66b1ff"}]}