# software/S9010/S9010-enterprise-layer.json
{"description": "Enterprise techniques used by GlassWorm, ATT&CK software S9010 (v1.0)", "name": "GlassWorm (S9010)", "domain": "enterprise-attack", "versions": {"layer": "4.5", "attack": "19", "navigator": "5.3.2"}, "techniques": [{"techniqueID": "T1071", "showSubtechniques": true}, {"techniqueID": "T1071.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has used HTTP for C2 and extracts data from the HTTP response headers.(Citation: Koi Glassworm InvisibleCode October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1560", "showSubtechniques": true}, {"techniqueID": "T1560.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has archived collected files within a zip file prior to exfiltration to include `/tmp/out.zip`.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1547", "showSubtechniques": true}, {"techniqueID": "T1547.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has set registry run keys for persistence in both `HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run` and `HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\`.(Citation: Koi Glassworm InvisibleCode October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1217", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has searched browser data for cookies, history, login databases, and cryptocurrency wallets.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1059", "showSubtechniques": true}, {"techniqueID": "T1059.002", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has utilized AppleScript to include `set keychainPassword to do shell script` to execute shell command that retrieves passwords from the macOS keychain.(Citation: Koi Glassworm New Tricks December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1059.007", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has leveraged JavaScript to execute its malicious code to include its hidden Unicode characters using the `eval` call.(Citation: Koi Glassworm Extensions November 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Aikido GlassWorm October 2025)(Citation: Socket GlassWorm January 2026)  [GlassWorm](https://attack.mitre.org/software/S9010) has also utilized encrypted payloads compiled in JavaScript.(Citation: Koi Glassworm New Tricks December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1554", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) can modify hardware wallet applications.(Citation: Koi Glassworm New Tricks December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1543", "showSubtechniques": true}, {"techniqueID": "T1543.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has established persistence on macOS via a LaunchAgent by writing a plist under `/library/LaunchAgents`.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1555", "showSubtechniques": true}, {"techniqueID": "T1555.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has collected keys stored within `/Library/Keychains/login.keychain-db`.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1555.003", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has gathered credentials stored in Mozilla FireFox and Chromium-based Browsers.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1602", "showSubtechniques": true}, {"techniqueID": "T1602.002", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has gathered data pertaining to VPN configurations.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Socket GlassWorm January 2026)  [GlassWorm](https://attack.mitre.org/software/S9010) has also targeted locally stored data on macOS located in `/Library/Application Support/Fortinet/FortiClient/conf/vpn.plist`.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1213", "showSubtechniques": true}, {"techniqueID": "T1213.003", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has gathered code repository authentication materials for NPM and GitHub.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Socket GlassWorm January 2026)  [GlassWorm](https://attack.mitre.org/software/S9010) has collected details pertaining to the npm configuration data for `_authToken`.(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1213.006", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has collected data from macOS devices through the gathering of Apple Notes related files by targeting `/Library/Group Containers/group.com.apple.notes/NoteStore.sqlite`,  `/Library/Group Containers/group.com.apple.notes/NoteStore.sqlite-wal`, and `/Library/Group Containers/group.com.apple.notes/NoteStore.sqlite-shm`.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1005", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has collected local data from a compromised host to include desktop cryptocurrency wallet data, and documents from within Desktop, Documents, and Downloads.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1565", "showSubtechniques": true}, {"techniqueID": "T1565.002", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) can intercept and modify transaction details associated with hardware wallet applications before signing.(Citation: Koi Glassworm New Tricks December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1074", "showSubtechniques": true}, {"techniqueID": "T1074.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has staged collected data in a working directory within a temp folder to include `/tmp/ijewf`.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1678", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has used a timeout function set to `9e5` which delays execution 900,000 milliseconds or 15 minutes to avoid detection.(Citation: Koi Glassworm New Tricks December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1140", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has decoded its Base64 instructions.(Citation: Koi Glassworm InvisibleCode October 2025)  [GlassWorm](https://attack.mitre.org/software/S9010) has also decrypted its AES protected payloads.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1480", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has utilized logic to avoid executing on Russian based devices.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1008", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has utilized Google Calendar as backup C2.(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Koi GlassWorm Rust December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1657", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has the ability to steal credentials for cryptocurrency wallets.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1564", "showSubtechniques": true}, {"techniqueID": "T1564.003", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has leveraged Hidden Virtual Network Computing (HVNC) to remain undetected and conduct execution of collection and communication actions.(Citation: Koi Glassworm InvisibleCode October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1105", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has downloaded additional payloads from C2.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm Extensions November 2025)(Citation: Socket GlassWorm January 2026)(Citation: Koi GlassWorm Rust December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1036", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has masqueraded as legitimate VSCode extensions.(Citation: Aikido GlassWorm October 2025)(Citation: Koi GlassWorm Rust December 2025)  [GlassWorm](https://attack.mitre.org/software/S9010) has also impersonated Github projects.(Citation: Aikido GlassWorm October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1571", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has distributed C2 using BitTorrent\u2019s Distributed Hash Table (DHT) network to harness a decentralized command capability.(Citation: Koi Glassworm InvisibleCode October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1027", "showSubtechniques": true}, {"techniqueID": "T1027.013", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has leveraged AES-256-CBC encryption to obfuscate its malicious JavaScript payload.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Socket GlassWorm January 2026)(Citation: Koi GlassWorm Rust December 2025)  [GlassWorm](https://attack.mitre.org/software/S9010) has also utilized Base64 encoding to obfuscate the C2 details stored in the Solana memo field.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Koi GlassWorm Rust December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1027.018", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has utilized invisible Unicode Private Use Area (PUA) characters to obfuscate its malicious code so that it does not render in code editors.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Aikido GlassWorm October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1090", "showSubtechniques": true}, {"techniqueID": "T1090.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has leveraged peer-to-peer software to facilitate communications within the victim network to include the software WebRTC.(Citation: Koi Glassworm InvisibleCode October 2025)  [GlassWorm](https://attack.mitre.org/software/S9010) has also established a SOCKS proxy to interact with victim devices that also acted as a proxy node for follow-on behaviors.(Citation: Koi Glassworm InvisibleCode October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1518", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has searched for existing wallet applications to include Ledger Live and Trezor Suite.(Citation: Koi Glassworm New Tricks December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1539", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has harvested Safari cookies stored within `/Library/Containers/com.apple.Safari/Data/Library/Cookies/ Cookies.binarycookies`.(Citation: Socket GlassWorm January 2026)  [GlassWorm](https://attack.mitre.org/software/S9010) has also stolen cookies within Chromium and Firefox browsers.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1195", "showSubtechniques": true}, {"techniqueID": "T1195.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has spread through Visual Studio extensions.(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Aikido GlassWorm October 2025)(Citation: Socket GlassWorm January 2026)  [GlassWorm](https://attack.mitre.org/software/S9010) has also spread through JavaScript projects hosted on Github.(Citation: Aikido GlassWorm October 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1082", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has the ability to check the OS of the victim host.(Citation: Socket GlassWorm January 2026)(Citation: Koi GlassWorm Rust December 2025)  [GlassWorm](https://attack.mitre.org/software/S9010) has checked whether the OS platform value includes `darwin` prior to execution of macOS specific scripts.(Citation: Socket GlassWorm January 2026)(Citation: Koi GlassWorm Rust December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1614", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has leveraged geofencing logic to detect whether it is operating in a Russian associated time zone to determine whether it continues to execute.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1614.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has identified the system language settings by checking for `ru_RU`, `ru-RU`, `ru`, and `Russian` to prevent execution in a Russian associated device.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1124", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has the ability to check the system\u2019s time zone on the victim device.(Citation: Socket GlassWorm January 2026)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1102", "showSubtechniques": true}, {"techniqueID": "T1102.001", "comment": "[GlassWorm](https://attack.mitre.org/software/S9010) has leveraged blockchain-based C2 infrastructure to include Solana blockchain that contains additional C2 details within the memo field.(Citation: Koi Glassworm New Tricks December 2025)(Citation: Koi Glassworm Extensions November 2025)(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Aikido GlassWorm October 2025)(Citation: Socket GlassWorm January 2026)(Citation: Koi GlassWorm Rust December 2025)  [GlassWorm](https://attack.mitre.org/software/S9010) has also leveraged Google Calendar to host encoded data.(Citation: Koi Glassworm InvisibleCode October 2025)(Citation: Socket GlassWorm January 2026)(Citation: Koi GlassWorm Rust December 2025)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by GlassWorm", "color": "#66b1ff"}]}