# versions/v7/software/S0148/S0148-enterprise-layer.json
{"description": "Enterprise techniques used by RTM, ATT&CK software S0148 v1.1", "name": "RTM (S0148)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"techniqueID": "T1548", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1548.002", "showSubtechniques": true, "comment": "RTM can attempt to run the program as admin, then show a fake error message and a legitimate UAC bypass prompt to the user in an attempt to socially engineer the user into escalating privileges.[1]"}, {"techniqueID": "T1071", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1071.001", "showSubtechniques": true, "comment": "RTM has initiated connections to external domains using HTTPS.[2] "}, {"score": 1, "techniqueID": "T1119", "showSubtechniques": false, "comment": "RTM monitors browsing activity and automatically captures screenshots if a victim browses to a URL matching one of a list of strings.[1][2]"}, {"techniqueID": "T1547", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1547.001", "showSubtechniques": true, "comment": "RTM tries to add a Registry Run key under the name \"Windows Update\" to establish persistence.[1]"}, {"score": 1, "techniqueID": "T1115", "showSubtechniques": false, "comment": "RTM collects data from the clipboard.[1][2]"}, {"techniqueID": "T1059", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1059.003", "showSubtechniques": true, "comment": "RTM uses the command line and rundll32.exe to execute.[1]"}, {"score": 1, "techniqueID": "T1568", "showSubtechniques": false, "comment": "RTM has resolved Pony C2 server IP addresses by either converting Bitcoin blockchain transaction data to specific octets, or accessing IP addresses directly within the Namecoin blockchain.[3][2]"}, {"techniqueID": "T1573", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1573.001", "showSubtechniques": true, "comment": "RTM encrypts C2 traffic with a custom RC4 variant.[1]"}, {"score": 1, "techniqueID": "T1083", "showSubtechniques": false, "comment": "RTM can check for specific files and directories associated with virtualization and malware analysis.[2]"}, {"score": 1, "techniqueID": "T1070", "showSubtechniques": true, "comment": "RTM has the ability to remove Registry entries that it created during execution.[1]"}, {"score": 1, "techniqueID": "T1070.004", "showSubtechniques": true, "comment": "RTM can delete all files created during its execution.[1][2]"}, {"score": 1, "techniqueID": "T1105", "showSubtechniques": false, "comment": "RTM can download additional files.[1][2]"}, {"techniqueID": "T1056", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1056.001", "showSubtechniques": true, "comment": "RTM can record keystrokes from both the keyboard and virtual keyboard.[1][2]"}, {"techniqueID": "T1559", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1559.002", "showSubtechniques": true, "comment": "RTM can search for specific strings within browser tabs using a Dynamic Data Exchange mechanism.[1]"}, {"score": 1, "techniqueID": "T1036", "showSubtechniques": true, "comment": "RTM has been delivered as archived Windows executable files masquerading as PDF documents.[2] "}, {"score": 1, "techniqueID": "T1036.004", "showSubtechniques": true, "comment": "RTM has named the scheduled task it creates \"Windows Update\".[2]"}, {"score": 1, "techniqueID": "T1112", "showSubtechniques": false, "comment": "RTM can delete all Registry entries created during its execution.[1]"}, {"score": 1, "techniqueID": "T1106", "showSubtechniques": false, "comment": "RTM can use the FindNextUrlCacheEntryA and FindFirstUrlCacheEntryA functions to search for specific strings within browser history.[1]"}, {"score": 1, "techniqueID": "T1571", "showSubtechniques": false, "comment": "RTM used Port 44443 for its VNC module.[1]"}, {"score": 1, "techniqueID": "T1027", "showSubtechniques": false, "comment": "RTM strings, network data, configuration, and modules are encrypted with a modified RC4 algorithm. RTM has also been delivered to targets as various archive files including ZIP, 7-ZIP, and RAR.[1][2]"}, {"score": 1, "techniqueID": "T1120", "showSubtechniques": false, "comment": "RTM can obtain a list of smart card readers attached to the victim.[1][2]"}, {"techniqueID": "T1566", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1566.001", "showSubtechniques": true, "comment": "RTM has been delivered via spearphishing attachments disguised as PDF documents.[2]"}, {"score": 1, "techniqueID": "T1057", "showSubtechniques": false, "comment": "RTM can obtain information about process integrity levels.[1]"}, {"score": 1, "techniqueID": "T1219", "showSubtechniques": false, "comment": "RTM has the capability to download a VNC module from command and control (C2).[1]"}, {"techniqueID": "T1053", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1053.005", "showSubtechniques": true, "comment": "RTM tries to add a scheduled task to establish persistence.[1][2]"}, {"score": 1, "techniqueID": "T1113", "showSubtechniques": false, "comment": "RTM can capture screenshots.[1][2]"}, {"techniqueID": "T1218", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1218.011", "showSubtechniques": true, "comment": "RTM runs its core DLL file using rundll32.exe.[1][2]"}, {"score": 1, "techniqueID": "T1518", "showSubtechniques": true, "comment": "RTM can scan victim drives to look for specific banking software on the machine to determine next actions.[1]"}, {"score": 1, "techniqueID": "T1518.001", "showSubtechniques": true, "comment": "RTM can obtain information about security software on the victim.[1]"}, {"techniqueID": "T1553", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1553.002", "showSubtechniques": true, "comment": "RTM samples have been signed with a code-signing certificates.[1]"}, {"score": 1, "techniqueID": "T1553.004", "showSubtechniques": true, "comment": "RTM can add a certificate to the Windows store.[1][2]"}, {"score": 1, "techniqueID": "T1082", "showSubtechniques": false, "comment": "RTM can obtain the computer name, OS version, and default language identifier.[1]"}, {"score": 1, "techniqueID": "T1033", "showSubtechniques": false, "comment": "RTM can obtain the victim username and permissions.[1]"}, {"score": 1, "techniqueID": "T1124", "showSubtechniques": false, "comment": "RTM can obtain the victim time zone.[1]"}, {"techniqueID": "T1204", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1204.002", "showSubtechniques": true, "comment": "RTM has relied on users opening malicious email attachments, decompressing the attached archive, and double-clicking the executable within.[2]"}, {"score": 1, "techniqueID": "T1497", "showSubtechniques": false, "comment": "RTM can detect if it is running within a sandbox or other virtualized analysis environment.[2]    "}, {"techniqueID": "T1102", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1102.001", "showSubtechniques": true, "comment": "RTM has used an RSS feed on Livejournal to update a list of encrypted C2 server names. RTM has also hidden Pony C2 server IP addresses within transactions on the Bitcoin and Namecoin blockchain.[1][3][2]"}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by RTM", "color": "#66b1ff"}]}