# versions/v5/groups/G0088/G0088-enterprise-layer.json
{"description": "Enterprise techniques used by TEMP.Veles, ATT&CK group G0088 v1.0", "name": "TEMP.Veles (G0088)", "domain": "mitre-enterprise", "version": "2.1", "techniques": [{"score": 1, "techniqueID": "T1043", "comment": "TEMP.Veles has used port 443 for C2."}, {"score": 1, "techniqueID": "T1003", "comment": "TEMP.Veles has used Mimikatz and a custom tool, SecHack, to harvest credentials."}, {"score": 1, "techniqueID": "T1074", "comment": "TEMP.Veles has created staging folders in directories that were infrequently used by legitimate users or processes."}, {"score": 1, "techniqueID": "T1133", "comment": "TEMP.Veles has used a VPN to persist in the victim environment."}, {"score": 1, "techniqueID": "T1107", "comment": "TEMP.Veles routinely deleted tools, logs, and other files after they were finished with them."}, {"score": 1, "techniqueID": "T1183", "comment": "TEMP.Veles has modified and added entries within HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options to maintain persistence.    "}, {"score": 1, "techniqueID": "T1066", "comment": "TEMP.Veles has modified files based on the open-source project cryptcat in an apparent attempt to decrease AV detection rates."}, {"score": 1, "techniqueID": "T1036", "comment": "TEMP.Veles has performed a variety of methods to look like valid users, including renaming files and mimicking legitimate administrator activities."}, {"score": 1, "techniqueID": "T1086", "comment": "TEMP.Veles has used a publicly-available PowerShell-based tool, WMImplant. The group has also used PowerShell to perform Timestomping."}, {"score": 1, "techniqueID": "T1076", "comment": "TEMP.Veles utilized RDP throughout an operation.   "}, {"score": 1, "techniqueID": "T1021", "comment": "TEMP.Veles has relied on encrypted SSH-based tunnels to transfer tools and for remote command/program execution."}, {"score": 1, "techniqueID": "T1053", "comment": "TEMP.Veles has used scheduled task XML triggers."}, {"score": 1, "techniqueID": "T1099", "comment": "TEMP.Veles used timestomping to modify the $STANDARD_INFORMATION attribute on tools."}, {"score": 1, "techniqueID": "T1065", "comment": "TEMP.Veles has used ports 4444, 8531, and 50501 for C2."}, {"score": 1, "techniqueID": "T1078", "comment": "TEMP.Veles has used compromised VPN accounts.  "}, {"score": 1, "techniqueID": "T1100", "comment": "TEMP.Veles has planted webshells on Outlook Exchange servers."}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by TEMP.Veles", "color": "#66b1ff"}]}