# versions/v11/software/S0504/S0504-enterprise-layer.json
{"description": "Enterprise techniques used by Anchor, ATT&CK software S0504 v1.0", "name": "Anchor (S0504)", "domain": "enterprise-attack", "versions": {"layer": "4.3", "attack": "11", "navigator": "4.5"}, "techniques": [{"techniqueID": "T1071", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1071.001", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) has used HTTP and HTTPS in C2 communications.(Citation: Cyberreason Anchor December 2019)"}, {"score": 1, "techniqueID": "T1071.004", "showSubtechniques": true, "comment": "Variants of [Anchor](https://attack.mitre.org/software/S0504) can use DNS tunneling to communicate with C2.(Citation: Cyberreason Anchor December 2019)(Citation: Medium Anchor DNS July 2020)"}, {"techniqueID": "T1059", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1059.003", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) has used cmd.exe to run its self deletion routine.(Citation: Cyberreason Anchor December 2019)"}, {"score": 1, "techniqueID": "T1059.004", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can execute payloads via shell scripting.(Citation: Medium Anchor DNS July 2020)"}, {"techniqueID": "T1543", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1543.003", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can establish persistence by creating a service.(Citation: Cyberreason Anchor December 2019)\t"}, {"score": 1, "techniqueID": "T1480", "showSubtechniques": false, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can terminate itself if specific execution flags are not present.(Citation: Cyberreason Anchor December 2019)"}, {"score": 1, "techniqueID": "T1008", "showSubtechniques": false, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can use secondary C2 servers for communication after establishing connectivity and relaying victim information to primary C2 servers.(Citation: Cyberreason Anchor December 2019)"}, {"techniqueID": "T1564", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1564.004", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) has used NTFS to hide files.(Citation: Cyberreason Anchor December 2019)"}, {"techniqueID": "T1070", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1070.004", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can self delete its dropper after the malware is successfully deployed.(Citation: Cyberreason Anchor December 2019)"}, {"score": 1, "techniqueID": "T1105", "showSubtechniques": false, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can download additional payloads.(Citation: Cyberreason Anchor December 2019)(Citation: Medium Anchor DNS July 2020)"}, {"score": 1, "techniqueID": "T1095", "showSubtechniques": false, "comment": "[Anchor](https://attack.mitre.org/software/S0504) has used ICMP in C2 communications.(Citation: Cyberreason Anchor December 2019)"}, {"score": 1, "techniqueID": "T1027", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) has obfuscated code with stack strings and string encryption.(Citation: Cyberreason Anchor December 2019)"}, {"score": 1, "techniqueID": "T1027.002", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) has come with a packed payload.(Citation: Cyberreason Anchor December 2019)"}, {"techniqueID": "T1021", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1021.002", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can support windows execution via SMB shares.(Citation: Medium Anchor DNS July 2020)"}, {"techniqueID": "T1053", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1053.003", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can install itself as a cron job.(Citation: Medium Anchor DNS July 2020)"}, {"score": 1, "techniqueID": "T1053.005", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can create a scheduled task for persistence.(Citation: Cyberreason Anchor December 2019)"}, {"techniqueID": "T1553", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1553.002", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) has been signed with valid certificates to evade detection by security tools.(Citation: Cyberreason Anchor December 2019)"}, {"score": 1, "techniqueID": "T1082", "showSubtechniques": false, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can determine the hostname and linux version on a compromised host.(Citation: Medium Anchor DNS July 2020)"}, {"score": 1, "techniqueID": "T1016", "showSubtechniques": false, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can determine the public IP and location of a compromised host.(Citation: Medium Anchor DNS July 2020)"}, {"techniqueID": "T1569", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1569.002", "showSubtechniques": true, "comment": "[Anchor](https://attack.mitre.org/software/S0504) can create and execute services to load its payload.(Citation: Cyberreason Anchor December 2019)(Citation: Medium Anchor DNS July 2020)"}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by Anchor", "color": "#66b1ff"}]}