# versions/v5/groups/G0087/G0087-enterprise-layer.json
{"description": "Enterprise techniques used by APT39, ATT&CK group G0087 v2.0", "name": "APT39 (G0087)", "domain": "mitre-enterprise", "version": "2.1", "techniques": [{"score": 1, "techniqueID": "T1090", "comment": "APT39 used custom tools to create SOCK5 proxies between infected hosts."}, {"score": 1, "techniqueID": "T1003", "comment": "APT39 has used Mimikatz, Ncrack, Windows Credential Editor and ProcDump to dump credentials."}, {"score": 1, "techniqueID": "T1002", "comment": "APT39 has used WinRAR and 7-Zip to compress an archive stolen data."}, {"score": 1, "techniqueID": "T1046", "comment": "APT39 used a custom port scanner known as BLUETORCH"}, {"score": 1, "techniqueID": "T1060", "comment": "APT39 has maintained persistence using the startup folder."}, {"score": 1, "techniqueID": "T1076", "comment": "APT39 has been seen using RDP for lateral movement and persistence."}, {"score": 1, "techniqueID": "T1021", "comment": "APT39 used secure shell (SSH) to move laterally among their targets."}, {"score": 1, "techniqueID": "T1053", "comment": "APT39 has created scheduled tasks."}, {"score": 1, "techniqueID": "T1064", "comment": "APT39 utilized custom scripts to perform internal reconnaissance."}, {"score": 1, "techniqueID": "T1023", "comment": "APT39 has modified LNK shortcuts."}, {"score": 1, "techniqueID": "T1045", "comment": "APT39 has repacked a modified version of Mimikatz to thwart anti-virus detection."}, {"score": 1, "techniqueID": "T1193", "comment": "APT39 leveraged spearphishing emails with malicious attachments to initially compromise victims."}, {"score": 1, "techniqueID": "T1192", "comment": "APT39 leveraged spearphishing emails with malicious links to initially compromise victims."}, {"score": 1, "techniqueID": "T1016", "comment": "APT39 has used nbtscan to discover vulnerable systems."}, {"score": 1, "techniqueID": "T1033", "comment": "APT39 used Remexi to collect usernames from the system."}, {"score": 1, "techniqueID": "T1204", "comment": "APT39 has sent spearphishing emails in an attempt to lure users to click on a malicious attachment or link."}, {"score": 1, "techniqueID": "T1078", "comment": "APT39 has used stolen credentials to compromise Outlook Web Access (OWA)."}, {"score": 1, "techniqueID": "T1100", "comment": "APT39 has installed ANTAK and ASPXSPY web shells."}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by APT39", "color": "#66b1ff"}]}