# versions/v7/groups/G0100/G0100-enterprise-layer.json
{"description": "Enterprise techniques used by Inception, ATT&CK group G0100 v1.0", "name": "Inception (G0100)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"techniqueID": "T1071", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1071.001", "showSubtechniques": true, "comment": "Inception has used HTTP, HTTPS, and WebDav in network communications.[3][1]"}, {"techniqueID": "T1547", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1547.001", "showSubtechniques": true, "comment": "Inception has maintained persistence by modifying Registry run key value  HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\.[3]"}, {"techniqueID": "T1059", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1059.005", "showSubtechniques": true, "comment": "Inception has used VBScript to execute malicious commands and payloads.[1][3]"}, {"score": 1, "techniqueID": "T1059.001", "showSubtechniques": true, "comment": "Inception has used PowerShell to execute malicious commands and payloads.[1][3]"}, {"techniqueID": "T1555", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1555.003", "showSubtechniques": true, "comment": "Inception used a browser plugin to steal passwords and sessions from Internet Explorer, Chrome, Opera, Firefox, Torch, and Yandex.[2]"}, {"score": 1, "techniqueID": "T1005", "showSubtechniques": false, "comment": "Inception used a file hunting plugin to collect .txt, .pdf, .xls or .doc files from the infected host.[4]"}, {"techniqueID": "T1573", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1573.001", "showSubtechniques": true, "comment": "Inception has encrypted network communications with AES.[3]"}, {"score": 1, "techniqueID": "T1203", "showSubtechniques": false, "comment": "Inception has exploited CVE-2012-0158, CVE-2014-1761, CVE-2017-11882 and CVE-2018-0802 for execution.[4][3][2][1]"}, {"score": 1, "techniqueID": "T1083", "showSubtechniques": false, "comment": "Inception used a file listing plugin to collect information about file and directories both on local and remote drives.[2]"}, {"score": 1, "techniqueID": "T1027", "showSubtechniques": false, "comment": "Inception has encrypted malware payloads dropped on victim machines with AES and RC4 encryption.[3]"}, {"techniqueID": "T1069", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1069.002", "showSubtechniques": true, "comment": "Inception has used specific malware modules to gather domain membership.[2]"}, {"techniqueID": "T1566", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1566.001", "showSubtechniques": true, "comment": "Inception has used weaponized documents attached to spearphishing emails for reconnaissance and initial compromise.[3][2][1][4]"}, {"score": 1, "techniqueID": "T1057", "showSubtechniques": false, "comment": "Inception has used a reconnaissance module to identify active processes and other associated loaded modules.[2]"}, {"techniqueID": "T1090", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1090.003", "showSubtechniques": true, "comment": "Inception used chains of compromised routers to proxy C2 communications between them and cloud service providers.[2]"}, {"techniqueID": "T1218", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1218.010", "showSubtechniques": true, "comment": "Inception has ensured persistence at system boot by setting the value regsvr32 %path%\\ctfmonrn.dll /s.[3]"}, {"score": 1, "techniqueID": "T1218.005", "showSubtechniques": true, "comment": "Inception has used malicious HTA files to drop and execute malware.[4]"}, {"score": 1, "techniqueID": "T1518", "showSubtechniques": false, "comment": "Inception has enumerated installed software on compromised systems.[2]"}, {"score": 1, "techniqueID": "T1082", "showSubtechniques": false, "comment": "Inception has used a reconnaissance module to gather information about the operating system and hardware on the infected host.[2]"}, {"score": 1, "techniqueID": "T1221", "showSubtechniques": false, "comment": "Inception has used decoy documents to load malicious remote payloads via HTTP.[1]"}, {"techniqueID": "T1204", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1204.002", "showSubtechniques": true, "comment": "Inception lured victims into clicking malicious files for machine reconnaissance and to execute malware.[3][4][2][1]"}, {"score": 1, "techniqueID": "T1102", "showSubtechniques": false, "comment": "Inception has incorporated at least five different cloud service providers into their C2 infrastructure including CloudMe.[3][2]"}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by Inception", "color": "#66b1ff"}]}