# software/S0512/S0512-enterprise-layer.json
{"description": "Enterprise techniques used by FatDuke, ATT&CK software S0512 (v1.1)", "name": "FatDuke (S0512)", "domain": "enterprise-attack", "versions": {"layer": "4.5", "attack": "17", "navigator": "5.1.0"}, "techniques": [{"techniqueID": "T1071", "showSubtechniques": true}, {"techniqueID": "T1071.001", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can be controlled via a custom C2 protocol over HTTP.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1547", "showSubtechniques": true}, {"techniqueID": "T1547.001", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) has used HKLM\\SOFTWARE\\Microsoft\\CurrentVersion\\Run to establish persistence.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1059", "showSubtechniques": true}, {"techniqueID": "T1059.001", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) has the ability to execute PowerShell scripts.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1005", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can copy files and directories from a compromised host.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1140", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can decrypt AES encrypted C2 communications.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1573", "showSubtechniques": true}, {"techniqueID": "T1573.001", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can AES encrypt C2 communications.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1008", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) has used several C2 servers per targeted organization.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1083", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can enumerate directories on target machines.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1070", "showSubtechniques": true}, {"techniqueID": "T1070.004", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can secure delete its DLL.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1036", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) has attempted to mimic a compromised user's traffic by using the same user agent as the installed browser.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1106", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can call ShellExecuteW to open the default browser on the URL localhost.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1027", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can use base64 encoding, string stacking, and opaque predicates for obfuscation.(Citation: ESET Dukes October 2019)", "score": 1, "showSubtechniques": true}, {"techniqueID": "T1027.002", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) has been regularly repacked by its operators to create large binaries and evade detection.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1027.016", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) has been packed with junk code and strings.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1057", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can list running processes on the localhost.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1090", "showSubtechniques": true}, {"techniqueID": "T1090.001", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can used pipes to connect machines with restricted internet access to remote machines via other infected hosts.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1012", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can get user agent strings for the default browser from HKCU\\Software\\Classes\\http\\shell\\open\\command.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1218", "showSubtechniques": true}, {"techniqueID": "T1218.011", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can execute via rundll32.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}, {"techniqueID": "T1082", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can collect the user name, Windows version, computer name, and available space on discs from a compromised host.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1016", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can identify the MAC address on the target computer.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": false}, {"techniqueID": "T1497", "showSubtechniques": true}, {"techniqueID": "T1497.003", "comment": "[FatDuke](https://attack.mitre.org/software/S0512) can turn itself on or off at random intervals.(Citation: ESET Dukes October 2019)", "score": 1, "color": "#66b1ff", "showSubtechniques": true}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by FatDuke", "color": "#66b1ff"}]}