# previous/july2019/software/S0223/S0223-enterprise-layer.json
{"description": "Enterprise techniques used by POWERSTATS, ATT&CK software S0223 v2.0", "name": "POWERSTATS (S0223)", "domain": "mitre-enterprise", "version": "2.1", "techniques": [{"score": 1, "techniqueID": "T1087", "comment": "POWERSTATS can retrieve usernames from compromised hosts."}, {"score": 1, "techniqueID": "T1043", "comment": "POWERSTATS has used port 80 for C2."}, {"score": 1, "techniqueID": "T1090", "comment": "POWERSTATS has connected to C2 servers through proxies."}, {"score": 1, "techniqueID": "T1132", "comment": "POWERSTATS encoded C2 traffic with base64."}, {"score": 1, "techniqueID": "T1005", "comment": "POWERSTATS can upload files from compromised hosts."}, {"score": 1, "techniqueID": "T1140", "comment": "POWERSTATS can deobfuscate the main backdoor code."}, {"score": 1, "techniqueID": "T1089", "comment": "POWERSTATS can disable Microsoft Office Protected View by changing Registry keys."}, {"score": 1, "techniqueID": "T1175", "comment": "POWERSTATS can use DCOM (targeting the 127.0.0.1 loopback address) to execute additional payloads on compromised hosts."}, {"score": 1, "techniqueID": "T1173", "comment": "POWERSTATS can use DDE to execute additional payloads on compromised hosts."}, {"score": 1, "techniqueID": "T1107", "comment": "POWERSTATS can delete all files on the C:\\, D:\\, E:\\ and, F:\\ drives using PowerShell Remove-Item commands."}, {"score": 1, "techniqueID": "T1036", "comment": "POWERSTATS has created a scheduled task named \"MicrosoftEdge\" to establish persistence."}, {"score": 1, "techniqueID": "T1170", "comment": "POWERSTATS can use Mshta.exe to execute additional payloads on compromised hosts."}, {"score": 1, "techniqueID": "T1027", "comment": "POWERSTATS uses character replacement, PowerShell environment variables, and XOR encoding to obfuscate code. POWERSTATS's backdoor code is a multi-layer obfuscated, encoded, and compressed blob."}, {"score": 1, "techniqueID": "T1086", "comment": "POWERSTATS uses PowerShell for obfuscation and execution."}, {"score": 1, "techniqueID": "T1105", "comment": "POWERSTATS can retrieve and execute additional PowerShell payloads from the C2 server."}, {"score": 1, "techniqueID": "T1053", "comment": "POWERSTATS has established persistence through a scheduled task using the command \u201dC:\\Windows\\system32\\schtasks.exe\u201d /Create /F /SC DAILY /ST 12:00 /TN MicrosoftEdge /TR \u201cc:\\Windows\\system32\\wscript.exe C:\\Windows\\temp\\Windows.vbe\u201d."}, {"score": 1, "techniqueID": "T1029", "comment": "POWERSTATS can sleep for a given number of seconds."}, {"score": 1, "techniqueID": "T1113", "comment": "POWERSTATS can retrieve screenshots from compromised hosts."}, {"score": 1, "techniqueID": "T1064", "comment": "POWERSTATS can use VBScript (VBE), PowerShell, and JavaScript code for execution."}, {"score": 1, "techniqueID": "T1063", "comment": "POWERSTATS has detected security tools."}, {"score": 1, "techniqueID": "T1032", "comment": "POWERSTATS has encrypted C2 traffic with RSA."}, {"score": 1, "techniqueID": "T1082", "comment": "POWERSTATS can retrieve OS name/architecture and computer/domain name information from compromised hosts."}, {"score": 1, "techniqueID": "T1016", "comment": "POWERSTATS can retrieve IP and network adapter configuration information from compromised hosts."}, {"score": 1, "techniqueID": "T1065", "comment": "POWERSTATS has used ports 8060 and 8888 for C2."}, {"score": 1, "techniqueID": "T1047", "comment": "POWERSTATS can use WMI queries to retrieve data from compromised hosts."}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by POWERSTATS", "color": "#66b1ff"}]}