# software/S0248/S0248-enterprise-layer.json
{"description": "Enterprise techniques used by yty, ATT&CK software S0248 v1.1", "name": "yty (S0248)", "domain": "mitre-enterprise", "version": "2.1", "techniques": [{"score": 1, "techniqueID": "T1009", "comment": "yty contains junk code in its binary, likely to confuse malware analysts."}, {"score": 1, "techniqueID": "T1005", "comment": "yty collects files with the following extensions: .ppt, .pptx, .pdf, .doc, .docx, .xls, .xlsx, .docm, .rtf, .inp, .xlsm, .csv, .odt, .pps, .vcf and sends them back to the C2 server."}, {"score": 1, "techniqueID": "T1083", "comment": "yty gathers information on victim\u2019s drives and has a plugin for document listing."}, {"score": 1, "techniqueID": "T1056", "comment": "yty uses a keylogger plugin to gather keystrokes."}, {"score": 1, "techniqueID": "T1036", "comment": "yty contains several references to football (including \"football,\" \"score,\" \"ball,\" and \"loose\") in a likely attempt to disguise its traffic."}, {"score": 1, "techniqueID": "T1057", "comment": "yty gets an output of running processes using the tasklist command."}, {"score": 1, "techniqueID": "T1018", "comment": "yty uses the net view command for discovery."}, {"score": 1, "techniqueID": "T1053", "comment": "yty establishes persistence by creating a scheduled task with the command SchTasks /Create /SC DAILY /TN BigData /TR \u201c + path_file + \u201c/ST 09:30\u201c."}, {"score": 1, "techniqueID": "T1113", "comment": "yty collects screenshots of the victim machine."}, {"score": 1, "techniqueID": "T1045", "comment": "yty packs a plugin with UPX."}, {"score": 1, "techniqueID": "T1082", "comment": "yty gathers the computer name, the serial number of the main disk volume, CPU information, Microsoft Windows version, and runs the command systeminfo."}, {"score": 1, "techniqueID": "T1016", "comment": "yty runs ipconfig /all and collects the domain name."}, {"score": 1, "techniqueID": "T1033", "comment": "yty collects the victim\u2019s username."}, {"score": 1, "techniqueID": "T1497", "comment": "yty has some basic anti-sandbox detection that tries to detect Virtual PC, Sandboxie, and VMware."}, {"score": 1, "techniqueID": "T1102", "comment": "yty communicates to the C2 server by retrieving a Google Doc."}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by yty", "color": "#66b1ff"}]}