# groups/G0101/G0101-enterprise-layer.json
{"description": "Enterprise techniques used by Frankenstein, ATT&CK group G0101 v1.0", "name": "Frankenstein (G0101)", "domain": "mitre-enterprise", "version": "3.0", "techniques": [{"score": 1, "techniqueID": "T1119", "showSubtechniques": false, "comment": "Frankenstein has enumerated hosts via Empire, gathering the username, domain name, machine name, and other system information.[1]"}, {"score": 1, "techniqueID": "T1020", "showSubtechniques": false, "comment": "Frankenstein has collected information via Empire, which is automatically sent the data back to the adversary's C2.[1]"}, {"techniqueID": "T1059", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1059.001", "showSubtechniques": true, "comment": "Frankenstein has used PowerShell to run a series of base64-encoded commands, that acted as a stager and enumerated hosts.[1]   "}, {"score": 1, "techniqueID": "T1059.003", "showSubtechniques": true, "comment": "Frankenstein has run a command script to set up persistence as a scheduled task named \"WinUpdate\", as well as other encoded commands from the command-line.[1]"}, {"score": 1, "techniqueID": "T1059.005", "showSubtechniques": true, "comment": "Frankenstein has used Word documents that prompts the victim to enable macros and run a Visual Basic script.[1]"}, {"score": 1, "techniqueID": "T1005", "showSubtechniques": false, "comment": "Frankenstein has enumerated hosts via Empire, gathering various local system information.[1]"}, {"score": 1, "techniqueID": "T1140", "showSubtechniques": false, "comment": "Frankenstein has deobfuscated base64-encoded commands following the execution of a malicious script, which revealed a small script designed to obtain an additional payload.[1]  "}, {"techniqueID": "T1573", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1573.001", "showSubtechniques": true, "comment": "Frankenstein has communicated with a C2 via an encrypted RC4 byte stream and AES-CBC.[1]"}, {"score": 1, "techniqueID": "T1041", "showSubtechniques": false, "comment": "Frankenstein has collected information via Empire, which is automatically sent the data back to the adversary's C2.[1]"}, {"score": 1, "techniqueID": "T1203", "showSubtechniques": false, "comment": "Frankenstein has used CVE-2017-11882 to execute code on the victim's machine.[1]"}, {"score": 1, "techniqueID": "T1105", "showSubtechniques": false, "comment": "Frankenstein has uploaded and downloaded files to utilize additional plugins.[1]"}, {"score": 1, "techniqueID": "T1027", "showSubtechniques": false, "comment": "Frankenstein has run encoded commands from the command line.[1]"}, {"score": 1, "techniqueID": "T1003", "showSubtechniques": false, "comment": "Frankenstein has harvested credentials from the victim's machine using Empire.[1]"}, {"techniqueID": "T1566", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1566.001", "showSubtechniques": true, "comment": "Frankenstein has used spearphishing emails to send trojanized Microsoft Word documents.[1]  "}, {"score": 1, "techniqueID": "T1057", "showSubtechniques": false, "comment": "Frankenstein has enumerated hosts, looking to obtain a list of all currently running processes.[1]"}, {"techniqueID": "T1053", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1053.005", "showSubtechniques": true, "comment": "Frankenstein has established persistence through a scheduled task using the command:  /Create /F /SC DAILY /ST 09:00 /TN WinUpdate /TR , named \"WinUpdate\".[1]    "}, {"techniqueID": "T1518", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1518.001", "showSubtechniques": true, "comment": "Frankenstein has used WMI queries to detect if virtualization environments or analysis tools were running on the system.[1]"}, {"score": 1, "techniqueID": "T1082", "showSubtechniques": false, "comment": "Frankenstein has enumerated hosts, looking for the system's machine name.[1]"}, {"score": 1, "techniqueID": "T1016", "showSubtechniques": false, "comment": "Frankenstein has enumerated hosts, looking for the public IP address of the system.[1]"}, {"score": 1, "techniqueID": "T1033", "showSubtechniques": false, "comment": "Frankenstein has enumerated hosts, gathering username, machine name, and administrative permissions information.[1]"}, {"score": 1, "techniqueID": "T1221", "showSubtechniques": false, "comment": "Frankenstein has used trojanized documents that retrieve remote templates from an adversary-controlled website.[1]"}, {"techniqueID": "T1127", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1127.001", "showSubtechniques": true, "comment": "Frankenstein has used MSbuild to execute an actor-created file.[1]"}, {"techniqueID": "T1204", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1204.002", "showSubtechniques": true, "comment": "Frankenstein has used trojanized Microsoft Word documents sent via email, which prompted the victim to enable macros.[1]"}, {"techniqueID": "T1497", "showSubtechniques": true}, {"score": 1, "techniqueID": "T1497.001", "showSubtechniques": true, "comment": "Frankenstein has used WMI queries to check if various security applications were running, including VMWare and Virtualbox.[1]"}, {"score": 1, "techniqueID": "T1047", "showSubtechniques": false, "comment": "Frankenstein has used WMI queries to check if various security applications were running, as well as the operating system version.[1]"}], "gradient": {"colors": ["#ffffff", "#66b1ff"], "minValue": 0, "maxValue": 1}, "legendItems": [{"label": "used by Frankenstein", "color": "#66b1ff"}]}